# templates/rtcwake.timer
[Unit]
Description=Run rtcwake everyday

[Timer]
OnCalendar=*-*-* {{ rtcwaketimer.hour }}
Persistent=true

[Install]
WantedBy=timers.target
